# debian/didiwiki.preinst
#! /bin/sh
# preinst script for didiwiki

set -e

case "$1" in
    install)
    ;;

    upgrade)
        if [ -x "/etc/init.d/didiwiki" ]; then
            if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
                    invoke-rc.d didiwiki stop || exit $?
            else   
                    /etc/init.d/didiwiki stop || exit $?
            fi
         fi
    ;;

    abort-upgrade)
    ;;

*)
    echo "preinst called with unknown argument '$1'" >&2
    exit 1
    ;;
esac

#DEBHELPER#

exit 0

